# Shell_Scripts/w/m/cms_imap_dir_replication/repl.cron
*/10 * * * * /root/lee/mail_repl/bin/mail_repl.sh > /root/lee/mail_repl/log/mail_repl.out 2> /root/lee/mail_repl/log/mail_repl.err
